# app/gumbode/Build/Products/Release/GUmBoDE.app/Contents/Resources/services/org.gumbode.privoxy.privoxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>gumbode.privoxy</string>
    <key>EnableGlobbing</key>
    <true/>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-c</string>
      <string>-l</string>
      <string>/usr/local/sbin/privoxy --no-daemon ~/.gumbode/privoxy/conf/config</string>
    </array>
    <key>StandardOutPath</key>
    <string>/tmp/privoxy.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/privoxy.log</string>
    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <false/>
    </dict>
  </dict>
</plist>
